# deploy/planka-mcp.service
[Unit]
Description=Planka MCP Server
After=network.target

[Service]
Type=simple
User=root
# Adjust this path to where your repository or compiled binary lives
WorkingDirectory=/opt/planka-mcp

# Environment Variables: Configure these securely using `systemctl edit planka-mcp`
# Environment=PLANKA_URL=...
# Environment=PLANKA_API_KEY=...

# Execute binary. Update this path to wherever your binary is installed (e.g., /usr/local/bin)
ExecStart=/usr/local/bin/planka-mcp --http --port 7526
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
